# dmr-nexus-sim.service
[Unit]
Description=DMR Nexus Simulated Repeater
After=dmr-nexus.service
Wants=dmr-nexus.service

[Service]
User=dmr-nexus
Group=dmr-nexus
WorkingDirectory=/opt/dmr-nexus
EnvironmentFile=/etc/default/dmr-nexus-sim
ExecStart=/opt/dmr-nexus/venv/bin/python /opt/dmr-nexus/tools/test_repeater.py \
  -s 127.0.0.1 \
  --base-id ${BASE_ID} \
  --callsign-prefix ${CALLSIGN} \
  --no-traffic
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
