# hblink4.service
[Unit]
Description=DMR Nexus Server
After=network.target

[Service]
Type=simple
User=dmr-nexus
Group=dmr-nexus
WorkingDirectory=/opt/dmr-nexus
ExecStart=/opt/dmr-nexus/venv/bin/python /opt/dmr-nexus/run.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Runtime directory for management socket
RuntimeDirectory=dmr-nexus

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmr-nexus

[Install]
WantedBy=multi-user.target
