# hblink4-dash.service
[Unit]
Description=DMR Nexus Dashboard
After=network.target dmr-nexus.service
Wants=dmr-nexus.service

[Service]
Type=simple
User=dmr-nexus
Group=dmr-nexus
WorkingDirectory=/opt/dmr-nexus/dashboard
ExecStart=/opt/dmr-nexus/venv/bin/python /opt/dmr-nexus/run_dashboard.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmr-nexus-dash

[Install]
WantedBy=multi-user.target
